# liquidsoap.initd.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          liquidsoap-daemon
# Required-Start:    $remote_fs $network $time
# Required-Stop:     $remote_fs $network $time
# Should-Start:
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Starts the liquidsoap daemon
# Description:
### END INIT INFO

liquidsoap=@liquidsoap_binary@
run_script=@run_script@
pid_dir=@pid_dir@
pid=`cat "${pid_dir}/run.pid" 2>/dev/null`

case "$1" in
  stop)
    if [ -z "${pid}" ]; then
      echo "No liquidsoap daemon running..";
      exit 0;
    fi;
    echo -n "Stopping liquidsoap... "
    kill -s TERM "${pid}"
    echo "OK"
    ;;

  start)
    if [ -n "${pid}" ]; then
      echo "PID file still present! Remove it if daemon isn't running..";
      exit 0;
    fi;
    echo -n "Starting liquidsoap... "
    ${liquidsoap} ${run_script}
    echo "OK"
    ;;

  restart|force-reload)
    $0 stop
    sleep 1
    $0 start
    ;;

  *)
    echo "Usage: $0 {start|stop|restart|force-reload}"
    exit 1
    ;;
esac
